# deploy/supervisor/supervisors/apiserver.ini
[program:apiserver]
environment=PYTHONPATH=%(ENV_PROJECT_ROOT_PATH)s:%(ENV_PROJECT_ROOT_PATH)s/app
command=%(ENV_PROJECT_ROOT_PATH)s/deploy/venv/bin/gunicorn apiserver:app -c %(ENV_PROJECT_ROOT_PATH)s/deploy/conf/gunicorn_cfg.py
process_name=%(program_name)s
redirect_stderr=false
stdout_logfile=none
stderr_logfile=none

[program:emailworker]
environment=PYTHONPATH=%(ENV_PROJECT_ROOT_PATH)s:%(ENV_PROJECT_ROOT_PATH)s/app
command=%(ENV_PROJECT_ROOT_PATH)s/deploy/venv/bin/celery worker -A apis.message.handlers.celery -P gevent --loglevel=info -Q send_email --pidfile=%(ENV_PROJECT_ROOT_PATH)s/deploy/supervisor/pids/%(program_name)s.pid
process_name=%(program_name)s
redirect_stderr=true
stdout_logfile=%(ENV_PROJECT_ROOT_PATH)s/logs/%(program_name)s.log
stderr_logfile=%(ENV_PROJECT_ROOT_PATH)s/logs/%(program_name)s.log

[program:smsworker]
environment=PYTHONPATH=%(ENV_PROJECT_ROOT_PATH)s:%(ENV_PROJECT_ROOT_PATH)s/app
command=%(ENV_PROJECT_ROOT_PATH)s/deploy/venv/bin/celery worker -A apis.message.handlers.celery -P gevent --loglevel=info -Q send_sms --pidfile=%(ENV_PROJECT_ROOT_PATH)s/deploy/supervisor/pids/%(program_name)s.pid
process_name=%(program_name)s
redirect_stderr=true
stdout_logfile=%(ENV_PROJECT_ROOT_PATH)s/logs/%(program_name)s.log
stderr_logfile=%(ENV_PROJECT_ROOT_PATH)s/logs/%(program_name)s.log
